# caffrey-ops-backup.timer
[Unit]
Description=Caffrey Ops nightly backup (03:00 daily)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=caffrey-ops-backup.service

[Install]
WantedBy=timers.target
